# share/docker/etc/crontabs/nobody
*/1 * * * * /var/www/scatbay.com/_articles/share/cron/drop-expired-html-cache.sh
