# share/docker/etc/services.d/03-cron/run
#!/bin/sh -e

echo "[services.d] starting service: cron"

crond -f -L /var/log/crond.log
