# etc/crontabs/nobody
*/1 * * * * /var/www/scatbay.com/_articles.git/share/cron/drop-expired-html-cache.sh
