# etc/cron.d/articles
*/1 * * * * www-data /var/www/scatbay.com/.articles.git/share/cron/drop-expired-html-cache.sh
